# examples/supervisord/redis.conf
[program:redis]
autostart = true
autorestart = true
command = "redis-server"
redirect_stderr = true
stdout_logfile = /var/log/supervisor/redis.log
